# rootfs/etc/s6-overlay/s6-rc.d/fix-perms/up
/command/with-contenv /etc/s6-overlay/scripts/fix-perms.sh
